# etc/update-motd.d/99-mymotd-generator
#!/bin/bash

# --- Настройка стилей (Возврат к старым цветам) ---
setup_colors() {
    tcDkG="\033[01;30m"   # Серый (лого)
    tcW="\033[01;37m"     # Белый
    YW="\033[33m"         # Желтый
    GN="\033[1;92m"       # Зеленый
    BL="\033[36m"         # Голубой
    RD="\033[01;31m"      # Красный
    CL="\033[m"           # Сброс
    tcRESET="\033[0m"
    TAB="  "
    
    ICON_HOST="${TAB}🏠${TAB}${CL}"
    ICON_IP="${TAB}💡${TAB}${CL}"
    ICON_DEF="${TAB}⚙️${TAB}${CL}"
    ICON_DISK="${TAB}💾${TAB}${CL}"
    ICON_LOAD="${TAB}📈${TAB}${CL}"
}

# --- Сбор данных ---
VIRT=$(systemd-detect-virt)
IP=$(hostname -I | awk '{print $1}')
HOSTNAME_F=$(hostname -f)
HOSTNAME_C=$(hostname)

uptime_s=$(cat /proc/uptime | cut -f1 -d.)
upD=$((uptime_s/60/60/24))
upH=$((uptime_s/60/60%24))
upM=$((uptime_s/60%60))

MEM_USED=$(free | grep Mem | awk '{printf "%.1f", $3/$2 * 100.0}')
SWAP_USED=$(free | grep Swap | awk '{if ($2>0) printf "%.1f", $3/$2 * 100.0; else print "0"}')
DISK_PERC=$(df -h / | awk 'NR==2 {print $5}')
DISK_FREE=$(df -h / | awk 'NR==2 {print $4}')
L_AVG=$(cut -d' ' -f1-3 /proc/loadavg)

if [ -f "/etc/os-release" ]; then
    OS_NAME=$(grep ^NAME /etc/os-release | cut -d= -f2 | tr -d '"')
    OS_VERSION=$(grep ^VERSION_ID /etc/os-release | cut -d= -f2 | tr -d '"')
fi

# --- Логотипы ---
if [ "$VIRT" = "lxc" ]; then
    raw_logo=$(cat << 'EOF'
 .____     ____  ____________  
 |    |    \   \/  /\_   ___ \ 
 |    |     \      / /    \  \/ 
 |    |___  /     \ \      \____
 |_______ \/___/\  \ \______  /
Container\/      \_/ aWasik \/ 
EOF
)
elif [ "$VIRT" = "kvm" ]; then
    raw_logo=$(cat << 'EOF'
____   ____________    _________
\   \ /   /\______ \  /   _____/
 \   Y   /  |    |  \ \_____  \ 
  \     /   |    `   \/        \
   \___/   /_______  /_______  /
                   \/ aWasik \/ 
EOF
)
else
    raw_logo=$(cat << 'EOF'
______________    _______________
\______   \   \ /   /\_   _____/
 |     ___/\   Y   /  |    __)_ 
 |    |     \     /   |        \
 |____|      \___/   /_______  /
                      aWasik \/ 
EOF
)
fi

# --- Логика выравнивания ---
setup_colors
mapfile -t logo_lines <<< "$raw_logo"
max_len=0
for line in "${logo_lines[@]}"; do
    [ ${#line} -gt $max_len ] && max_len=${#line}
done

logo_width=$((max_len + 2))
label_w=16

print_line() {
    local idx=$1
    local icon=$2
    local l_color=$3
    local l_text=$4
    local val=$5
    
    local img_part="${logo_lines[$idx]}"
    
    # Теперь и текст, и значение используют l_color, как на твоем скрине
    printf "${tcDkG}%-${logo_width}s${icon} ${l_color}%-${label_w}s %s${CL}\n" \
        "$img_part" "$l_text" "$val"
}

# --- Финальный вывод ---
echo -e ""
print_line 0 "$ICON_HOST" "$GN" "Hostname:"      "$HOSTNAME_F ($HOSTNAME_C)"
print_line 1 "$ICON_IP"   "$GN" "IP Address:"    "$IP"
print_line 2 "$ICON_DEF"  "$BL" "${OS_NAME}:"    "$OS_VERSION"
print_line 3 "$ICON_DEF"  "$BL" "Load Average:"  "$L_AVG"
print_line 4 "$ICON_DEF"  "$BL" "Disk Usage:"    "$DISK_PERC ($DISK_FREE free)"
print_line 5 "$ICON_DEF"  "$BL" "Memory used:"   "${MEM_USED}% (Swap: ${SWAP_USED}%)"
print_line 6 "$ICON_DEF"  "$BL" "System uptime:" "${upD}d ${upH}h ${upM}m"
echo -e "${tcRESET}"

echo "🔑 Last login...........: $(last | head -1 | cut -c 1-9 | xargs) at $(last | head -1 | cut -c 40-55 | xargs) from $(last | head -1 | cut -c 23-39 | xargs)"
echo "⏳ Uptime...............: $(uptime -p | sed 's/up //')"
CPU_MODEL=$(grep -m1 "model name" /proc/cpuinfo | cut -d ':' -f2 | sed 's/^ //')
echo "🖥  CPU Model............: $CPU_MODEL"
echo "⚙  Processes............: $(ps -ef | wc -l) running processes"
echo "🧬 OS...................: $(lsb_release -ds 2>/dev/null || grep PRETTY_NAME /etc/os-release | cut -d= -f2 | tr -d '\"')"

NET_IFACE=$(ip route get 1.1.1.1 2>/dev/null | awk '{for(i=1;i<=NF;i++) if($i=="dev") print $(i+1)}')
if [ -n "$NET_IFACE" ]; then
      RX_BYTES=$(cat /sys/class/net/$NET_IFACE/statistics/rx_bytes)
      TX_BYTES=$(cat /sys/class/net/$NET_IFACE/statistics/tx_bytes)

      function human_readable {
        local BYTES=$1
        local UNITS=('B' 'KB' 'MB' 'GB' 'TB')
        local UNIT=0
        local VALUE=$BYTES

        while (( VALUE >= 1024 && UNIT < 4 )); do
          VALUE=$((VALUE / 1024))
          ((UNIT++))
        done

        local DIVISOR=$((1024 ** UNIT))
        local FLOAT_VAL=$(awk "BEGIN {printf \"%.2f\", $BYTES / $DIVISOR}")

        if (( UNIT >= 3 )); then
          echo "${FLOAT_VAL} ${UNITS[$UNIT]}"
        else
          echo "$((BYTES / DIVISOR)) ${UNITS[$UNIT]}"
        fi
      }

      RX_HR=$(human_readable $RX_BYTES)
      TX_HR=$(human_readable $TX_BYTES)
      echo "🌐 Network Traffic......: Received: $RX_HR | Transmitted: $TX_HR"
else
      echo "🌐 Network Traffic......: Interface not found"
fi

if command -v ufw &>/dev/null; then
      UFW_STATUS=$(ufw status | head -1)
      UFW_RULES=$(ufw status numbered | grep -c '\[')
      echo "🛡  Firewall (UFW).......: $UFW_STATUS, Rules: $UFW_RULES"
else
      echo "🛡  Firewall (UFW).......: not installed"
fi
if command -v docker &>/dev/null; then
    RUNNING_CONTAINERS=$(docker ps -q | wc -l)
    TOTAL_CONTAINERS=$(docker ps -a -q | wc -l)
    echo "🐳 Docker containers....: $RUNNING_CONTAINERS / $TOTAL_CONTAINERS running:"
    if [ "$RUNNING_CONTAINERS" -gt 0 ]; then
      docker ps --format "   • {{.Names}}" | while read line; do printf "   %s\n" "$line"; done
    fi
else
    echo "🐳 Docker...............: not installed"
fi
